# scripts/tailscale.service
[Unit]
Description=Pulsar alerts for Tailnet

[Service]
Type=simple

ExecStart=/pulsar-scripts/tailscale.sh
Restart=on-failure
RestartSec=5
